# stable-server.service
[Unit]
Description=Think AI Stable Server with Streaming and Qwen
After=network.target ollama.service

[Service]
Type=exec
User=administrator
WorkingDirectory=/home/administrator/think_ai
Environment="RUST_LOG=info"
ExecStart=/home/administrator/think_ai/target/release/stable-server-streaming
Restart=always
RestartSec=10
StandardOutput=append:/home/administrator/think_ai/stable-server.log
StandardError=append:/home/administrator/think_ai/stable-server.log

[Install]
WantedBy=multi-user.target