# think-ai.service
[Unit]
Description=Think AI Full System Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=champi
Group=champi
WorkingDirectory=/home/champi/Dev/think_ai
Environment=RUST_LOG=info
Environment=PORT=7777
Environment=DEEPGRAM_API_KEY=e31341c95ee93fd2c8fced1bf37636f042fe038b
Environment=ELEVENLABS_API_KEY=sk_aa595f99bc5842b5df837d4c47fe3a18ce00b9a39a0f2877
Environment=AUDIO_CACHE_DIR=/home/champi/Dev/think_ai/full-system/audio_cache
ExecStart=/home/champi/Dev/think_ai/target/release/think-ai-full
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
LimitCORE=infinity

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/champi/Dev/think_ai

[Install]
WantedBy=multi-user.target
EOF < /dev/null
